# packaging/debian/finguard.service
[Unit]
Description=FinGuard - WireGuard Proxy with Web Management
Documentation=https://github.com/JPKribs/FinGuard
After=network-online.target avahi-daemon.service
Wants=network-online.target
RequiresMountsFor=/var/lib/finguard

[Service]
Type=simple
User=finguard
Group=finguard
ExecStart=/usr/local/bin/finguard --config /etc/finguard/config.yaml
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
RestartPreventExitStatus=2

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/finguard /var/log/finguard /run/finguard
PrivateTmp=true
ProtectKernelTunables=false
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=false
RestrictNamespaces=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Network capabilities for TUN device creation
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=finguard

# Working directory
WorkingDirectory=/var/lib/finguard

[Install]
WantedBy=multi-user.target