# packaging/debian/postinst
#!/bin/bash
set -e

case "$1" in
    configure)
        # MARK: Force hostname to FinGuard with persistence
        echo "Setting hostname to FinGuard with persistence..."
        
        # Set current hostname
        hostnamectl set-hostname FinGuard
        
        # Ensure persistence across multiple methods
        echo "FinGuard" > /etc/hostname
        
        # Update /etc/hosts to prevent resolution issues
        sed -i '/127\.0\.1\.1/d' /etc/hosts 2>/dev/null || true
        echo "127.0.1.1    FinGuard" >> /etc/hosts
        
        # For systems that use /etc/machine-info
        if [ -f /etc/machine-info ]; then
            sed -i '/PRETTY_HOSTNAME/d' /etc/machine-info 2>/dev/null || true
            echo 'PRETTY_HOSTNAME="FinGuard"' >> /etc/machine-info
        fi
        
        # For Armbian/embedded systems that might have additional config
        if [ -f /etc/armbian-release ]; then
            # Some Armbian images reset hostname from this file
            if [ -f /boot/armbianEnv.txt ]; then
                sed -i '/hostname/d' /boot/armbianEnv.txt 2>/dev/null || true
                echo "hostname=FinGuard" >> /boot/armbianEnv.txt
            fi
        fi

        # MARK: Create finguard user and group
        if ! getent group finguard > /dev/null 2>&1; then
            echo "Creating finguard group..."
            if command -v addgroup >/dev/null 2>&1; then
                addgroup --system --quiet finguard
            else
                groupadd --system finguard
            fi
        fi

        if ! getent passwd finguard > /dev/null 2>&1; then
            echo "Creating finguard user..."
            if command -v adduser >/dev/null 2>&1; then
                adduser --system --quiet \
                    --home /var/lib/finguard \
                    --no-create-home \
                    --shell /usr/sbin/nologin \
                    --ingroup finguard \
                    --gecos "FinGuard System User" \
                    finguard
            else
                useradd --system \
                    --home-dir /var/lib/finguard \
                    --no-create-home \
                    --shell /usr/sbin/nologin \
                    --gid finguard \
                    --comment "FinGuard System User" \
                    finguard
            fi
        fi

        # MARK: Add user to necessary groups
        if getent group netdev >/dev/null 2>&1; then
            usermod -a -G netdev finguard 2>/dev/null || true
        fi

        # MARK: Create directories with proper ownership
        install -d -o finguard -g finguard -m 0755 /var/lib/finguard
        install -d -o finguard -g finguard -m 0755 /var/log/finguard
        
        if [ ! -d /run/finguard ]; then
            install -d -o finguard -g finguard -m 0755 /run/finguard
        fi

        # MARK: Create symlink for web UI
        echo "Creating web UI symlink..."
        if [ -d /usr/local/share/finguard/web ] && [ ! -e /var/lib/finguard/web ]; then
            ln -sf /usr/local/share/finguard/web /var/lib/finguard/web
            echo "Web UI symlink created"
        fi

        # MARK: Install Avahi if needed
        if ! command -v avahi-daemon >/dev/null 2>&1; then
            echo "Installing avahi-daemon..."
            apt update && apt install -y avahi-daemon avahi-utils
        fi
        
        echo "Avahi configuration installed via package"

        # MARK: Set proper ownership for config files
        if [ -f /etc/finguard/config.yaml ]; then
            chown root:finguard /etc/finguard/config.yaml
            chmod 640 /etc/finguard/config.yaml
        fi
        
        if [ -f /etc/finguard/wireguard.yaml ]; then
            chown root:finguard /etc/finguard/wireguard.yaml
            chmod 600 /etc/finguard/wireguard.yaml
        fi

        # MARK: Set capabilities on binary
        if [ -f /usr/local/bin/finguard ]; then
            echo "Setting network capabilities on FinGuard binary..."
            if command -v setcap >/dev/null 2>&1; then
                setcap cap_net_admin,cap_net_raw,cap_net_bind_service+ep /usr/local/bin/finguard || {
                    echo "Warning: Failed to set capabilities. You may need to run as root."
                }
            fi
        fi

        # MARK: Enable and start services
        systemctl daemon-reload
        systemctl enable avahi-daemon
        systemctl restart avahi-daemon
        systemctl enable finguard.service
        
        echo "FinGuard installation completed successfully!"
        echo "Hostname set to: FinGuard (with persistence)"
        echo ""
        echo "Next steps:"
        echo "1. Edit /etc/finguard/config.yaml and replace the admin token"
        echo "2. systemctl start finguard"
        echo "3. Access web interface at http://localhost:10000"
        
        ;;
esac

exit 0