# packaging/debian/postinst
#!/bin/bash
set -e

case "$1" in
    configure)
        REAL_BINARY="/usr/local/lib/finguard/bin/finguard"

        # MARK: Ensure system user and group exist
        if ! getent group finguard >/dev/null 2>&1; then
            echo "Creating finguard group..."
            groupadd --system finguard
        fi

        if ! getent passwd finguard >/dev/null 2>&1; then
            echo "Creating finguard user..."
            useradd --system \
                --home-dir /var/lib/finguard \
                --no-create-home \
                --shell /usr/sbin/nologin \
                --gid finguard \
                --comment "FinGuard System User" \
                finguard
        fi

        # MARK: Create application data directories
        echo "Creating directories..."
        install -d -o finguard -g finguard -m 0775 /var/lib/finguard
        install -d -o finguard -g finguard -m 0775 /var/log/finguard
        install -d -o finguard -g finguard -m 0775 /etc/finguard/backups
        install -d -o finguard -g finguard -m 0775 /var/lib/finguard/backups

        # MARK: Prepare binary directory (must be writable for self-updates)
        install -d -o finguard -g finguard -m 0775 /usr/local/lib/finguard
        install -d -o finguard -g finguard -m 0775 /usr/local/lib/finguard/bin

        # MARK: Secure binary and set capabilities if possible
        if [ -f "$REAL_BINARY" ]; then
            echo "Setting ownership on FinGuard binary..."
            chown finguard:finguard "$REAL_BINARY"
            chmod 755 "$REAL_BINARY"

            if command -v setcap >/dev/null 2>&1; then
                setcap cap_net_admin,cap_net_raw,cap_net_bind_service+ep "$REAL_BINARY" 2>/dev/null || \
                    echo "Note: Capabilities will be provided by systemd"
            fi
        fi

        # MARK: Secure config files
        echo "Setting up config files..."
        chown -R finguard:finguard /etc/finguard
        chmod 660 /etc/finguard/*.yaml 2>/dev/null || true

        # MARK: Set ownership for web assets
        if [ -d /usr/local/share/finguard/web ]; then
            chown -R finguard:finguard /usr/local/share/finguard
        fi

        # MARK: Set hostname
        echo "FinGuard" > /etc/hostname
        hostnamectl set-hostname FinGuard

        # MARK: Ensure /etc/hosts maps hostname
        sed -i "s/127\.0\.1\.1.*/127.0.1.1 FinGuard/" /etc/hosts || \
        echo "127.0.1.1 FinGuard" >> /etc/hosts

        # MARK: Register systemd service
        systemctl daemon-reload
        systemctl enable finguard.service

        echo "========================================="
        echo "FinGuard installation completed!"
        echo "========================================="
        echo ""
        echo "✓ Service user created"
        echo "✓ Permissions configured"
        echo "✓ Auto-updates ready (NO SUDO NEEDED!)"
        echo ""
        echo "Next steps:"
        echo "1. Edit /etc/finguard/config.yaml"
        echo "2. systemctl start finguard"
        echo "3. Access web UI at http://localhost:10000"
        echo "========================================="
        ;;
esac

exit 0