# roles/FinGuard/templates/finguard-update.timer.j2
[Unit]
Description=Run FinGuard update every Wednesday at 03:00

[Timer]
OnCalendar=Wed *-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target