# roles/FinGuard/templates/status-update.timer.j2
[Unit]
Description=Update system status information every 30 seconds

[Timer]
OnBootSec=30s
OnUnitActiveSec=30s
AccuracySec=100ms
Unit=status-update.service

[Install]
WantedBy=timers.target